# app/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

# ── Nginx web server ────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/nginx-err.log
stdout_logfile=/var/log/supervisor/nginx-out.log
priority=10

# ── PHP-FPM ─────────────────────────────────────────────
[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/phpfpm-err.log
stdout_logfile=/var/log/supervisor/phpfpm-out.log
priority=5
